# crontab/crontab
0 0 * * * find /var/lib/php/sessions -atime +7 | xargs -r rm
